# templates/default/carbon.debian.init.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          carbon-<%= @name %>-<%= @instance %>
# Required-Start:    $network $named $remote_fs $syslog
# Required-Stop:     $network $named $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO 
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=<%= node['graphite']['carbon']['cache_bin'] %>
NAME="carbon-<%= @name %>-<%= @instance %>"
PIDFILE="<%= node['graphite']['carbon']['pid_dir'] %>/${NAME}.pid"
DESC="Graphite cache daemon, instance <%= @instance %>"

test -x $DAEMON || exit 0

[ ! -d <%= node['graphite']['carbon']['pid_dir'] %> ] && mkdir -p <%= node['graphite']['carbon']['pid_dir'] %>
chown <%= node['graphite']['user_account'] %> <%= node['graphite']['carbon']['pid_dir'] %>

set -e

case "$1" in
  start)
  echo -n "Starting $DESC: "
  start-stop-daemon --start --quiet \
    --exec ${DAEMON} --pidfile ${PIDFILE} -- --config=<%= node['graphite']['carbon']['conf_dir'] %>/carbon.conf --pid=${PIDFILE} \
    --instance=<%= @instance %> start
  echo "$NAME."
  ;;
  stop)
  echo -n "Stopping $DESC: "
        start-stop-daemon --stop  --signal KILL --pidfile ${PIDFILE} \
                2>&1 > /dev/null
  echo "$NAME."
  ;;
  status)
        start-stop-daemon --status --quiet --pidfile ${PIDFILE}
  ;;
  reload)
  ;;
  restart|force-reload)
  $0 stop
  $0 start
  ;;
  *)
  N=/etc/init.d/$NAME
  # echo "Usage: $N {start|stop|restart|reload|force-reload}" >&2
  echo "Usage: $N {start|stop|restart|force-reload}" >&2
  exit 1
  ;;
esac

exit 0

